# billing_server/deploy/frogswork-auto-billing.timer
[Unit]
Description=Daily FrogsWork auto-billing timer

[Timer]
OnCalendar=*-*-* 06:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
